# deploy/netmessenger-web.service
[Unit]
Description=NETMessenger Web (Flask + Gunicorn)
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/netmessenger-web
EnvironmentFile=/etc/netmessenger/web.env
ExecStart=/var/www/netmessenger-web/venv/bin/gunicorn -c deploy/gunicorn.conf.py
Restart=on-failure
RestartSec=5

# Безопасность
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
